# supervisord.conf
[supervisord] 
nodaemon=true 
[program:mycat] 
command=/bin/bash /root/mycat/bin/mycat start 